# helper/supervisord.conf
; vasili-helper supervisord configuration.
;
; Each program is a launcher shim that reads /etc/vasili-helper/config.json
; and either exec's the real daemon or sleeps forever if the section is
; disabled. The helper UI restarts the affected programs via `supervisorctl`
; whenever the operator clicks Apply.

[supervisord]
nodaemon=true
logfile=/var/log/vasili-helper/supervisord.log
pidfile=/run/supervisord.pid
user=root

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:helper-ui]
command=/usr/bin/python3 -u /opt/vasili-helper/app/app.py
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/vasili-helper/helper-ui.log
stderr_logfile=/var/log/vasili-helper/helper-ui.log

[program:dns-proxy]
command=/opt/vasili-helper/launchers/dns-proxy.sh
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/vasili-helper/dns-proxy.log
stderr_logfile=/var/log/vasili-helper/dns-proxy.log

[program:crack-server]
command=/opt/vasili-helper/launchers/crack-server.sh
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/vasili-helper/crack-server.log
stderr_logfile=/var/log/vasili-helper/crack-server.log

[program:iodine-backend]
command=/opt/vasili-helper/launchers/iodine-backend.sh
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/vasili-helper/iodine-backend.log
stderr_logfile=/var/log/vasili-helper/iodine-backend.log

[program:wg-backend]
command=/opt/vasili-helper/launchers/wg-backend.sh
autostart=true
autorestart=true
priority=45
stdout_logfile=/var/log/vasili-helper/wg-backend.log
stderr_logfile=/var/log/vasili-helper/wg-backend.log

[program:ssh-tunnel]
command=/opt/vasili-helper/launchers/ssh-tunnel.sh
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/vasili-helper/ssh-tunnel.log
stderr_logfile=/var/log/vasili-helper/ssh-tunnel.log
